# drifter-repo/services/drifter-dashboard.service
[Unit]
Description=MZ1312 DRIFTER Web Dashboard & Audio Bridge
After=mosquitto.service
Wants=mosquitto.service

StartLimitIntervalSec=120
StartLimitBurst=5
[Service]
Type=simple
WorkingDirectory=/opt/drifter
ExecStart=/opt/drifter/venv/bin/python3 /opt/drifter/web_dashboard.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
