# drifter-repo/services/drifter-thresholds.service
[Unit]
Description=MZ1312 DRIFTER Adaptive Thresholds
After=nanomq.service
Wants=nanomq.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/opt/drifter
ExecStart=/opt/drifter/venv/bin/python3 /opt/drifter/adaptive_thresholds.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/opt/drifter/.env

[Install]
WantedBy=multi-user.target
