# drifter-repo/services/drifter-watchdog.service
[Unit]
Description=MZ1312 DRIFTER Watchdog — Service & System Health Monitor
After=nanomq.service drifter-canbridge.service
Wants=nanomq.service

[Service]
Type=simple
ExecStart=/opt/drifter/venv/bin/python3 /opt/drifter/watchdog.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
