# services/drifter-alerts.service
[Unit]
Description=MZ1312 DRIFTER Alert Engine
After=nanomq.service drifter-canbridge.service
Wants=nanomq.service
# Vehicle-survival hardening 2026-06-10: NO StartLimitAction=reboot-force.
# This consumes CAN telemetry (After=drifter-canbridge); the old
# reboot-force + StartLimitBurst=5 meant a crash here — e.g. on a no-CAN /
# K-line car, or a broker hiccup at boot — could reboot the whole node, the
# exact reboot-loop "never brick" forbids. A down alert engine is already
# surfaced by /healthz (it's not hw-optional → degraded/503), so the operator
# still sees it; recovery is Restart=on-failure (below) + the watchdog
# restarting the unit — never a node reboot.

[Service]
User=drifter
Group=drifter
Type=simple
WorkingDirectory=/opt/drifter
ExecStart=/opt/drifter/venv/bin/python3 /opt/drifter/alert_engine.py
Restart=on-failure
RestartSec=15
Environment=PYTHONUNBUFFERED=1

# Vehicle-survival hardening 2026-05-20: driver-safety service —
# bias global OOM killer away from us. See progress.txt for the
# Ollama-uncapped OOM-priority-inversion discussion.
OOMScoreAdjust=-500

# --- Hardening (see AGENTS.md) --------------------------------------
# Baseline sandboxing. ProtectSystem=full makes /usr, /boot, /etc
# read-only while leaving /opt/drifter writable for logs/state.
NoNewPrivileges=true
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectSystem=full
LockPersonality=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
