# services/drifter-aidiag.service
[Unit]
Description=MZ1312 DRIFTER AI Diagnostics (Tier 2 Claude)
After=nanomq.service
Wants=nanomq.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/opt/drifter
ExecStart=/opt/drifter/venv/bin/python3 /opt/drifter/ai_diagnostics.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/opt/drifter/.env

# RAM-survival hardening 2026-06-10: heavy/sacrificial tier. LLM calls
# go out-of-process (Claude API / ollama), so own RSS is modest — 512M
# covers prompt/response buffers without ever covering the LLM runtime.
# OOMScoreAdjust=+500 makes the kernel kill this BEFORE the diag/safety
# core (at -800). See drifter-canbridge.service for the scheme.
MemoryMax=512M
OOMScoreAdjust=500

[Install]
WantedBy=multi-user.target
