# services/drifter-autoconnect.service
[Unit]
Description=MZ1312 DRIFTER Wi-Fi Hotspot Auto-Connector
After=NetworkManager.service nanomq.service
Wants=NetworkManager.service
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/opt/drifter
# nmcli connection changes need root (NetworkManager polkit). Configure the
# phone hotspot via /opt/drifter/.env: PHONE_HOTSPOT_SSID / PHONE_HOTSPOT_PSK.
ExecStart=/opt/drifter/venv/bin/python3 /opt/drifter/auto_connect.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/opt/drifter/.env

# RAM-survival hardening 2026-06-10: diag/safety core tier. Negative
# OOMScoreAdjust protects this from global OOM kill so the kernel
# sacrifices the heavy AI/ML services (+500) first. See
# drifter-canbridge.service for the full OOM-priority scheme.
OOMScoreAdjust=-800

[Install]
WantedBy=multi-user.target
