# services/drifter-location.service
[Unit]
Description=MZ1312 DRIFTER Location Service (Google Elevation + Places)
After=nanomq.service network-online.target
Wants=nanomq.service network-online.target
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/opt/drifter
ExecStart=/opt/drifter/venv/bin/python3 /opt/drifter/location_service.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/opt/drifter/.env

# RAM-survival hardening 2026-06-10: diag/safety core tier. Negative
# OOMScoreAdjust protects this from global OOM kill so the kernel
# sacrifices the heavy AI/ML services (+500) first. See
# drifter-canbridge.service for the full OOM-priority scheme.
OOMScoreAdjust=-800

[Install]
WantedBy=multi-user.target
