# services/drifter-discord.service
[Unit]
Description=MZ1312 DRIFTER Vivi Discord Bot
After=nanomq.service network.target
Wants=nanomq.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/opt/drifter
ExecStart=/opt/drifter/venv/bin/python3 /opt/drifter/vivi_discord.py
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/opt/drifter/.env

[Install]
WantedBy=multi-user.target
