# services/drifter-llm.service
[Unit]
Description=MZ1312 DRIFTER LLM Mechanic (superseded — disabled)
After=network.target drifter-canbridge.service
Wants=network.target

[Service]
User=drifter
Group=drifter
Type=simple
WorkingDirectory=/opt/drifter
Environment=PYTHONPATH=/opt/drifter
Environment=OLLAMA_HOST=localhost
Environment=OLLAMA_PORT=11434
ExecStart=/opt/drifter/venv/bin/python3 /opt/drifter/llm_mechanic.py
Restart=always
RestartSec=10
KillMode=process
TimeoutStopSec=10

NoNewPrivileges=true
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectSystem=full
LockPersonality=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
