# services/drifter-vision.service
[Unit]
Description=MZ1312 DRIFTER Vision Engine (YOLO/Hailo)
After=nanomq.service
Wants=nanomq.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/opt/drifter
ExecStart=/opt/drifter/venv/bin/python3 /opt/drifter/vision_engine.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/opt/drifter/.env

# RAM-survival hardening 2026-06-10: heavy/sacrificial tier. YOLO model +
# frame buffers are RAM-hungry in-process (Hailo offloads inference but
# the Python side still buffers frames/tensors, ~0.5-1.5G); 1.5G caps a
# runaway with headroom. OOMScoreAdjust=+500 makes the kernel kill this
# BEFORE the diag/safety core (at -800). See drifter-canbridge.service
# for the OOM-priority scheme.
MemoryMax=1.5G
OOMScoreAdjust=500

[Install]
WantedBy=multi-user.target
