# services/drifter-canbridge.service
[Unit]
Description=MZ1312 DRIFTER CAN Bridge
After=network.target nanomq.service
Wants=nanomq.service
# NO StartLimitAction=reboot-force here (removed 2026-06-10).
# A missing telemetry source — no USB2CANFD plugged in, no K-line adapter,
# adapter popped off mid-drive, or the ECU simply being off — is a NORMAL,
# expected condition on ANY vehicle. It must DEGRADE, never reboot the node.
# The old reboot-force + StartLimitBurst=3 + RestartSec=15 combination
# rebooted the Pi roughly every ~45s, forever, whenever CAN was absent —
# a self-inflicted reboot loop that bricks the node on a bench or on a car
# whose adapter dropped out. can_bridge.py now stays alive, retries CAN
# interface detection in-process, and publishes a 'hw-pending'/'no-can'
# status instead of exiting — so this unit should essentially never fail.
# A plain Restart=on-failure remains below purely as a last-resort safety net
# (real crash / OOM), with NO StartLimitAction so it can never escalate to a
# reboot. StartLimitBurst/StartLimitIntervalSec are dropped along with it.

[Service]
Type=simple
WorkingDirectory=/opt/drifter
ExecStartPre=/usr/local/bin/drifter-setup-can
ExecStart=/opt/drifter/venv/bin/python3 /opt/drifter/can_bridge.py
Restart=on-failure
RestartSec=15
Environment=PYTHONUNBUFFERED=1

# Vehicle-survival hardening 2026-05-20: canbridge has bounded ring
# buffers — 256M ceiling, growth past this is a coding defect.
#
# === RAM-survival OOM-priority scheme (2026-06-10) ===
# When the 8G Pi runs out of memory the kernel OOM-killer must sacrifice
# the heavy AI/ML/voice services FIRST and keep this diagnostics + safety
# core alive. We bias OOMScoreAdjust in two tiers:
#   * diagnostics/safety CORE (DIAG_SERVICES in config.py)  -> -800
#     (this file is core; strengthened from the historical -500)
#   * heavy/sacrificial (vivi, analyst, reporter, voicein,
#     fly-catcher, aidiag, vision)                          -> +500
# Lower score = last to be killed; higher = first. Ollama is
# intentionally LEFT UNCAPPED / unscored (user decision; see
# progress.txt OOM-priority-inversion discussion) — the +500 heavy tier
# is what the kernel reaches for before ever touching this core.
# canbridge is the primary telemetry source, so it gets full protection.
MemoryMax=256M
OOMScoreAdjust=-800

# --- Hardening (see AGENTS.md) --------------------------------------
# Baseline sandboxing. ProtectSystem=full makes /usr, /boot, /etc
# read-only while leaving /opt/drifter writable for logs/state.
NoNewPrivileges=true
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectSystem=full
LockPersonality=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
