# services/drifter-fly-catcher.service
[Unit]
Description=MZ1312 DRIFTER Fly Catcher — ADS-B ghost detection
After=nanomq.service
Wants=nanomq.service
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/opt/drifter
ExecStart=/opt/drifter/venv/bin/python3 /opt/drifter/fly_catcher.py
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

# RAM-survival hardening 2026-06-10: heavy/sacrificial tier. Loads an
# in-process ML model (~0.5-1G resident); 1.5G caps a runaway with
# headroom over real need. OOMScoreAdjust=+500 makes the kernel kill
# this BEFORE the diag/safety core (at -800). See
# drifter-canbridge.service for the OOM-priority scheme.
MemoryMax=1.5G
OOMScoreAdjust=500

NoNewPrivileges=true
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectSystem=full
LockPersonality=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
