# services/drifter-rf.service
[Unit]
Description=MZ1312 DRIFTER RF Monitor — RTL-SDR TPMS & Spectrum Scanner
After=nanomq.service
Wants=nanomq.service
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/opt/drifter
ExecStart=/opt/drifter/venv/bin/python3 /opt/drifter/rf_monitor.py
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
